# etc/cron/crontab
# m h  dom mon dow   command
0 5 1 */2 * /home/treescope/treescope-core/etc/cron/cron_job.sh
